# service-files/cerebro_api.service
[Unit]
Description=Cerebro API Service
After=network.target

[Service]
User=agentuser
Group=agentuser
WorkingDirectory=/app/cerebro/apps/server
ExecStart=/home/agentuser/.bun/bin/bun run start
Restart=always
RestartSec=2
Environment=NODE_ENV=production  # Add any environment variables you need

[Install]
WantedBy=multi-user.target
